# ansible/roles/photoprism/templates/rclone.service.j2
[Unit]
Description=Download latest pictures using rclone

[Service]
Type=simple
ExecStart={{ photoprism_base_path }}/rclone.sh
