# ansible/roles/photoprism/templates/rclone.timer.j2
[Unit]
Description=Run rclone sync every 6 hours

[Timer]
OnUnitActiveSec=6h
Persistent=true
Unit=rclone.service

[Install]
WantedBy=timers.target
