# ansible/roles/rclone/templates/rclone.service.j2
[Unit]
Description=Download latest pictures using rclone

[Service]
Type=simple
ExecStart={{ ansible_facts['user_dir'] }}/rclone.sh
