# ansible/roles/betisier/templates/betisier-reset.timer.j2
[Unit]
Description=Daily Betisier password reset
Requires=betisier-reset.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
